# systemd/travel-api.service
[Unit]
Description=Travel Marketplace Backend Docker Compose
Documentation=https://github.com/yourusername/travel-marketplace
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/travel-marketplace-backend
ExecStart=/usr/bin/docker compose -f /opt/travel-marketplace-backend/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /opt/travel-marketplace-backend/docker-compose.prod.yml down
ExecReload=/usr/bin/docker compose -f /opt/travel-marketplace-backend/docker-compose.prod.yml restart
TimeoutStartSec=0
Restart=on-failure
RestartSec=5

# Security settings
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/travel-marketplace-backend

[Install]
WantedBy=multi-user.target

